# image/patternjournal.cron
# Backup the db every day 7am UTC (midnight mountain time)
# The Backup script is on the image as part of the deployment.

0 7 * * * /patternjournal/image/backup_db.sh
